# deploy/systemd/dsp-client.service
[Unit]
Description=DarkSideProtocol Client
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/darksideprotocol
Environment=DSP_CLIENT_CONFIG=/opt/darksideprotocol/configs/client.json
ExecStart=/opt/darksideprotocol/bin/dsp-client
Restart=always
RestartSec=2
LimitNOFILE=1048576

[Install]
WantedBy=multi-user.target
